# deploy/systemd/deepcommodity@.service
[Unit]
Description=deepCommodity routine: %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=trader
Group=trader
WorkingDirectory=/srv/deepCommodity
EnvironmentFile=/srv/deepCommodity/.env
ExecStart=/srv/deepCommodity/deploy/run_routine.sh %i
# run_routine.sh owns per-routine logging (logs/%i.log) as the trader user.
# Do NOT add StandardOutput=append here — systemd opens append targets as root,
# which then blocks the trader process from writing the same file (Permission denied).
# systemd captures any stray stdout/stderr to the journal.

# Sandboxing — keep the routine away from anything outside the repo.
# /home/trader is writable so headless claude can persist its config/cache.
ProtectSystem=strict
ReadWritePaths=/srv/deepCommodity /home/trader
ProtectHome=read-only
NoNewPrivileges=true
PrivateTmp=true
LockPersonality=true
RestrictSUIDSGID=true
